# scripts/launchagents/com.helix.scheduled-agents.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.helix.scheduled-agents</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/jadinrobinson/ai-assistant/scripts/start_scheduled_agents.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/jadinrobinson/ai-assistant/backend</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
    <key>SCHEDULED_AGENTS_INTERVAL_SECONDS</key>
    <string>300</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/jadinrobinson/ai-assistant/backend/logs/scheduled-agents.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jadinrobinson/ai-assistant/backend/logs/scheduled-agents.err.log</string>
</dict>
</plist>
